# supervisor/newsletter-queue.conf
[program:newsletter_queue]
command=/usr/bin/php /path/to/wp-coupon-mailing/scripts/queue_coupons.php
autostart=false
autorestart=false
startsecs=0

