# supervisor/newsletter-send.conf
[program:newsletter_send]
command=/usr/bin/php /path/to/wp-coupon-mailing/scripts/send_coupons.php
autostart=false
autorestart=false
startsecs=0
